# config/camera.service
[Unit]
Description=DIY Camera Service
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u /opt/camera/app.py
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
