# config/camera.service
[Unit]
Description=DIY Camera Service
After=local-fs.target dev-fb1.device
Requires=dev-fb1.device

[Service]
Type=simple
ExecStart=/usr/bin/python3 -u /opt/camera/camera_app.py
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target